# scripts/sd/criu.service
[Unit]
Description=Checkpoint Restore in Userspace daemon

[Service]
ExecStart=/usr/sbin/criu service
